# supervisor/nginx.conf
[program:elasticsearch]
command=service nginx restart
