# supervisor/nginx.conf
[program:nginx]
command=service nginx restart
autorestart = true
